# post_update.sh
#!/bin/sh

# shellcheck disable=SC1091,2154
. /etc/rc.subr && load_rc_config zigbee2mqtt
